# _etc/systemd/system/cvesearch.web.service
[Unit]
Description=circl dot lu CVE-Search Web Server
Requires=mongodb.service
After=network.target cvesearch.db_init.service cvesearch.db_repopulate.service mongodb.service
Documentation=https://cve-search.github.io/cve-search/webgui/webgui.html

[Service]
WorkingDirectory=$(INSTALLDIR)
ExecStart=python3 ./web/index.py
User=$(CVE_USER)
Type=simple
SyslogIdentifier=cvesearch.web
Restart=always
RestartSec=5
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target
