# _etc/systemd/system/cvesearch.db_repopulate.target
[Unit]
Conflicts=cvesearch.web.service
Wants=cvesearch.db_repopulate.service
Before=cvesearch.db_repopulate.service
After=cvesearch.db_init.service
