# _etc/systemd/system/cvesearch.timer
[Unit]
Description=circl dot lu CVE search service trigger timer

[Timer]
Unit=cvesearch.service
OnCalendar=*-*-* *:0/30:00

[Install]
WantedBy=timers.target
