# _etc/systemd/system/cvesearch.db_init.service
[Unit]
Description=circl dot lu CVE-Search database initialization (only run this once!)
Requires=mongodb.service redis-server.service
After=network.target mongodb.service redis-server.service
RefuseManualStart=true
Documentation=https://cve-search.github.io/cve-search/getting_started/installation.html

[Service]
WorkingDirectory=$(INSTALLDIR)
ExecStart=python3 ./sbin/db_mgmt_cpe_dictionary.py -p
ExecStart=python3 ./sbin/db_mgmt_json.py -p
ExecStart=python3 ./sbin/db_updater.py -c
User=$(CVE_USER)
Type=oneshot
SyslogIdentifier=cvesearch.db_updater
