# _etc/systemd/system/cvesearch.db_repopulate.service
[Unit]
Description=circl dot lu CVE-Search database repopulation (drop & reimport everything)
Requires=mongod.service redis-server.service
After=network.target cvesearch.db_init.service cvesearch.db_updater.service mongod.service redis-server.service
RefuseManualStart=true
Documentation=https://cve-search.github.io/cve-search/database/database.html

[Service]
WorkingDirectory=/opt/cve/cve-search
ExecStart=python3 ./sbin/db_updater.py -f
User=cve
Type=oneshot
SyslogIdentifier=cvesearch.db_updater
