# _etc/systemd/system/cvesearch.db_repopulate.service
[Unit]
Description=circl dot lu CVE-Search database repopulation (drop & reimport everything)
Requires=mongodb.service redis-server.service
After=network.target cvesearch.db_init.service cvesearch.db_updater.service mongodb.service redis-server.service
RefuseManualStart=true
Documentation=https://cve-search.github.io/cve-search/database/database.html

[Service]
WorkingDirectory=$(INSTALLDIR)
ExecStart=python3 ./sbin/db_updater.py -v -f
User=$(CVE_USER)
Type=oneshot
SyslogIdentifier=cvesearch.db_updater
