# _etc/systemd/system/cvesearch.web.service
[Unit]
Description=circl dot lu CVE-Search Web Server
Requires=mongod.service
After=network.target cvesearch.db_init.service cvesearch.db_repopulate.service mongod.service
Documentation=https://cve-search.github.io/cve-search/webgui/webgui.html

[Service]
WorkingDirectory=/opt/cve/cve-search
ExecStart=python3 ./web/index.py
User=cve
Type=simple
SyslogIdentifier=cvesearch.web
Restart=always
RestartSec=5
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target
