# _etc/systemd/system/cvesearch.db_updater.timer
[Unit]
Description=circl dot lu CVE-Search db_updater trigger timer

[Timer]
Unit=cvesearch.db_updater.service
OnActiveSec=0
OnUnitActiveSec=7200
OnUnitInactiveSec=7200
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
