# _etc/systemd/system/cvesearch.db_init.target
[Unit]
Conflicts=cvesearch.web.service
Wants=cvesearch.db_init.service
Before=cvesearch.db_init.service
