# _etc/systemd/system/cvesearch.db_updater.timer
[Unit]
Description=circl dot lu CVE-Search db_updater trigger timer

[Timer]
Unit=cvesearch.db_updater.service
OnCalendar=*-*-* 5:15:00

[Install]
WantedBy=timers.target
